# laradock/workspace/crontab/laradock
*/5 * * * * laradock cd /var/www/; php bin/console app:cron >> /dev/null 2>&1
